# burp/files/burp.cron
#! /bin/sh

log=/var/log/burp

echo -n "====== " >> $log
date >> $log
/usr/sbin/burp -a t >> $log 2>&1
echo >> $log

