# roles/burp/files/burp.cron
#! /bin/sh

log=/var/log/burp

PATH=/usr/sbin:$PATH; export PATH

echo -n "====== " >> $log
date >> $log
burp -a t >> $log 2>&1
echo >> $log

